# scripts/mcp/systemd/mcp-watchdog.timer
[Unit]
Description=Run MCP watchdog every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true
Unit=mcp-watchdog.service

[Install]
WantedBy=timers.target
